# hermes-claude-proxy.service
[Unit]
Description=Claude CLI Subscription Adapter (Hermes Proxy)
Documentation=https://github.com/RawSmokeTerribilus/claudehermessubscriptionadapter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.hermes/proxy
ExecStart=%h/.hermes/proxy/venv/bin/python %h/.hermes/proxy/server.py --host 127.0.0.1 --port 8082
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-proxy

[Install]
WantedBy=default.target
